# service/mqtt-service/mqtt-service.service
[Unit]
Description=AIOT MQTT Service - 独立的MQTT消息处理服务
After=network.target mysql.service mosquitto.service
Wants=mysql.service mosquitto.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/codehubot/service/mqtt-service
Environment="PATH=/opt/codehubot/venv/bin:/usr/local/bin:/usr/bin:/bin"

# 从环境文件加载配置
EnvironmentFile=/opt/codehubot/service/mqtt-service/.env

# 启动命令
ExecStart=/opt/codehubot/venv/bin/python /opt/codehubot/service/mqtt-service/main.py

# 自动重启配置
Restart=always
RestartSec=5

# 优雅停止
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# 资源限制
MemoryLimit=512M
CPUQuota=50%

# 日志配置
StandardOutput=append:/var/log/codehubot/mqtt-service.log
StandardError=append:/var/log/codehubot/mqtt-service-error.log
SyslogIdentifier=mqtt-service

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target







